# BA/coffee/Thesis_Simon/coffee-v3/systemd/coffee-v3-download-comics.timer
[Unit]
Description=Run download-comics daily at 3am

[Timer]
OnCalendar=*-*-* 3:00:00

[Install]
WantedBy=timers.target